# Misc/Batch/PROJ_batch_loader.timer
[Unit]
Description=Batch downloader schedule

[Timer]
OnBootSec=1min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
